# s6/etc/services.d/nsssd/run
#!/bin/execlineb -P

# TODO convert to s6-rc format
# https://github.com/just-containers/s6-overlay/blob/master/MOVING-TO-V3.md
# https://skarnet.org/software/s6-rc/s6-rc-compile.html#source

fdmove -c 2 1
s6-envuidgid containizen
fdmove 1 3
s6-ipcserver -1 -U -- s
nsssd-unix
